# supervisor/workers.conf
[program:request_worker]
directory=/usr/src/app
command=celery --app=transcriptionservice.workers.celeryapp worker -n %(ENV_SERVICE_NAME)s_request_worker@%%h --queues=%(ENV_SERVICE_NAME)s_requests -c 4 --loglevel=INFO
priority=1

[program:ingress]
directory=/usr/src/app
command=python /usr/src/app/transcriptionservice/server/ingress.py --debug